# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command = /usr/bin/php-fpm -c /etc/php5/etc/php.ini -F
process_name = php-fpm
numprocs = 1
autostart = true
autorestart = true
user = root

